# ai.openclaw.wechat-ce-patch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.openclaw.wechat-ce-patch</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/w2ayz/.openclaw/workspace/skills/wechat-ce/boot-patch.sh</string>
    </array>

    <!-- Run once at login -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Also re-run whenever the plugin file is overwritten by openclaw update -->
    <key>WatchPaths</key>
    <array>
        <string>/Users/w2ayz/.openclaw/npm/node_modules/@tencent-weixin/openclaw-weixin/dist/src/messaging/process-message.js</string>
    </array>

    <key>StandardOutPath</key>
    <string>/tmp/openclaw/wechat-ce-patch.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/openclaw/wechat-ce-patch.log</string>
</dict>
</plist>
